# src/JT808.Servers/JT808.MsgId0x0200Server/JT808.MsgId0x0200Server.conf
[program:JT808.MsgId0x0200Server]
command=dotnet JT808.MsgId0x0200Server.dll
directory=/home/gps
user=root
stdout_logfile=/home/supervisor/logs/JT808MsgId0x0200Server.log
redirect_stderr=true