# src/GPS.Gateway.JT808SuperSocketServer/GPS.Gateway.JT808SuperSocketServer.conf
[program:GPS.Gateway.JT808SuperSocketServer]
command=mono GPS.Gateway.JT808SuperSocketServer.exe 
directory=/home/gps
user=root
stdout_logfile=/home/supervisor/logs/gps.log
redirect_stderr=true