# src/JT808.Netty/GPS.JT808NettyServer/GPS.JT808NettyServer.conf
[program:GPS.JT808NettyServer]
command=dotnet GPS.JT808NettyServer.dll --ASPNETCORE_ENVIRONMENT=Production
directory=/home/gps
user=root
stdout_logfile=/home/supervisor/logs/gps.log
redirect_stderr=true