# src/JT808.Service/JT808.MsgId0x0200Service/JT808.MsgId0x0200Service.conf
[program:JT808.MsgId0x0200Service]
command=dotnet JT808.MsgId0x0200Service.dll
directory=/home/gps
user=root
stdout_logfile=/home/supervisor/logs/MsgId0x0200Service.log
redirect_stderr=true