# src/GPS.BusinessServices/GPS.JT808SampleDeviceMonitoring/GPS.JT808SampleDeviceMonitoring.conf
[program:GPS.JT808SampleDeviceMonitoring]
command=dotnet GPS.JT808SampleDeviceMonitoring.dll
directory=/home/gps
user=root
stdout_logfile=/home/supervisor/logs/JT808SampleDeviceMonitoring.log
redirect_stderr=true